# cloudflare-updater.service
[Unit]
Description="Cloudflare record updater"

[Service]
ExecStart=./cloudflare.sh
